# packages/sysdef/sysdef-0.0.4-py3-none-any.whl/sysdef/res/sysdef.service
[Unit]
Description=sysdef
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/sysdef --boot

[Install]
WantedBy=multi-user.target
